# config/diana.ini
[supervisord]
nodaemon=true
logfile=/code/log/supervisord.log
pidfile=/code/log/supervisord.pid

[program:diana_web]
process_name = %(program_name)s-%(process_num)s
command=python app.py -debug=${DEBUG}
directory=/code
user=root
stopsignal=HUP
stopasgroup=true
killasgroup=true
autorestart=true

stdout_logfile=/code/log/diana_web_access_%(process_num)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

redirect_stderr=true
stderr_logfile=/code/log/diana_web_error_%(process_num)s.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5

numprocs = 1
numprocs_start = 8004

[group:diana]
programs=diana_web
